# conf/systemd/mumudvb@.service
[Unit]
Description=mumudvb %I
After=multi-user.target

[Service]
Type=simple
User=_mumudvb
Restart=always
RestartSec=10
ExecStart=/bin/bash -c '/home/pi/workspace/dvbstream/mumudvb.sh %I >> /home/pi/workspace/dvbstream/log/mumudvb-%I.log 2>&1'

[Install]
WantedBy=multi-user.target

